# net-mgmt/bandwidthd/pkg-plist
bin/bandwidthd
@sample etc/bandwidthd/bandwidthd.conf.sample
share/bandwidthd/htdocs/legend.gif
share/bandwidthd/htdocs/logo.gif
etc/rc.d/bandwidthd
%%PORTDOCS%%%%DOCSDIR%%/README
